# packages/pi-controller/systemd/lawnbot-controller.service
[Unit]
Description=LawnBot Controller
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=lds00
WorkingDirectory=/home/lds00/lawnbot
ExecStart=/usr/bin/python3 /home/lds00/lawnbot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
MemoryLimit=512M

[Install]
WantedBy=multi-user.target
